# supervisord.conf
[program:node]
command = node server.js
direcotry = /home/dotcloud/current
